# systemd/project-backup@.service
[Unit]
Description=Backup %i SQLite databases

[Service]
Type=oneshot
ExecStart=/usr/local/bin/project-backup %i
